# ec2ppuku.conf
description "ec2ppuku"

start on runlevel [2]
stop on runlevel [016]

script
  . /etc/profile
  while [ 1 ]
  do
    /usr/bin/ec2ppuku
    sleep 1
  done
end script
